# templates/prometheus.service.j2
[Unit]
Description=prometheus
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/opt/prometheus/prometheus
Restart=always

[Install]
WantedBy=default.target